# deploy/lucawriter.service
[Unit]
Description=LucaWriter SaaS service
After=network-online.target coobox.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/lucawriter-core
EnvironmentFile=/home/ubuntu/lucawriter.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/ubuntu/lucawriter-core/.venv/bin/python backend/main.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=45
KillSignal=SIGINT
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/home/ubuntu/lucawriter-data

[Install]
WantedBy=multi-user.target
